# deploy/supervisor/cihrms-queue-audit.conf
; CIHRMS — audit queue worker
; Place at /etc/supervisor/conf.d/cihrms-queue-audit.conf

[program:cihrms-queue-audit]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/cihrms/artisan queue:work --queue=audit --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/cihrms/queue-audit.log
stopwaitsecs=3600
