# deploy/supervisor/cihrms-scheduler.conf
; CIHRMS — Laravel scheduler (replaces cron entry)
; Place at /etc/supervisor/conf.d/cihrms-scheduler.conf

[program:cihrms-scheduler]
process_name=%(program_name)s
command=php /var/www/cihrms/artisan schedule:work
autostart=true
autorestart=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/cihrms/scheduler.log
