# deploy/supervisor/cihrms-queue-default.conf
; CIHRMS — default queue worker
; Place at /etc/supervisor/conf.d/cihrms-queue-default.conf
;
; Drains the `default` queue, which is where all Laravel notifications,
; the N1 SendSmsJob (SMS dispatcher), the N2 module listeners
; (loan / benefit / attendance / payroll / offboarding / asset /
; document notifications), and the Paystack webhook job all land.
;
; Without this worker running, those side-effects silently queue and
; never fire. The analytics queue worker (cihrms-queue-analytics.conf)
; does NOT cover the default queue.
;
; numprocs=2 is a sensible baseline for a single-host deploy. Scale up
; if PayrollRunPaid fan-out (one job per employee on a run) becomes a
; bottleneck — those jobs are short but a 500-employee run produces
; ~500 jobs in a burst.

[program:cihrms-queue-default]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/cihrms/artisan queue:work --queue=default --sleep=3 --tries=3 --max-time=3600 --backoff=10
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/log/cihrms/queue-default.log
stopwaitsecs=3600
